# scripts/livefeed.service
[Unit]
Description=Endoscopy Suite Camera Feed (GStreamer MJPEG to Display + TCP port 5000)
# ─────────────────────────────────────────────────────────────
# NOTE: This GStreamer service is NO LONGER REQUIRED for the
# camera pipeline. The pi_capture_daemon.js now captures
# directly from /dev/video0 using ffmpeg.
#
# To disable this service on the Pi:
#   sudo systemctl stop livefeed
#   sudo systemctl disable livefeed
#
# The daemon (managed by PM2) is the sole camera consumer now.
# GStreamer is only kept here for reference / display overlay.
# ─────────────────────────────────────────────────────────────
After=network.target
Wants=network.target

[Service]
Type=simple
User=lm
Group=lm
WorkingDirectory=/home/lm/loyalmed
ExecStart=/home/lm/loyalmed/scripts/livefeed_pi.sh
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

# Environment — camera parameters
Environment=VIDEO_DEVICE=/dev/video0
Environment=WIDTH=1280
Environment=HEIGHT=720
Environment=FRAMERATE=30

# Allow time for camera device + DRM to be ready
ExecStartPre=/bin/sleep 5

[Install]
WantedBy=multi-user.target
